# files/configs/profile.d/testnet_env_vars.sh.j2
#!/usr/bin/env sh

export INFLUXDB_CONNECTION_STRING='{{ vault_influxdb_connection_string }}'
